# systemd/nwperf.target
[Unit]
Description=NWPerf Services
After=syslog.target network.target

[Install]
WantedBy=multi-user.target
